# modules/base/root/etc/services.d/logging/run
#!/usr/bin bash

sleep 10

exec s6-setuidgid abc tail -F /usr/local/var/log/gvm/*